# os-setup/duckdns.service
[Unit]
Description=Update DuckDNS registration

[Service]
Type=simple
WorkingDirectory=/usr/share/duckdns
ExecStart=/bin/bash duckdns.sh